# rootfs.img.gz/fwo_rootfs.img/etc/upstart/lab126_gui_setup.conf
start on framework_ready and started kb and started pillow

task

env RESTART_JOB=lab126_gui

script
    source /etc/upstart/monitor_functions

    f_monitor_session_clear_count ${RESTART_JOB}
    f_monitor_persist_clear_count ${RESTART_JOB}
end script