# rootfs.img.gz/fwo_rootfs.img/etc/upstart/lab126.conf
# this is just a place holder to provide a collapse point for
# non-X and non-framework lab126 daemons
start on started system_setup_after_shpm and started demd
stop on stopping dbus or ota-update

export LANG LC_ALL
