# rootfs.img.gz/fwo_rootfs.img/etc/upstart/langpicker.conf
# start language picker after blanket has loaded the langpicker module

start on x_setup_ready and blanket_loaded_langpicker
stop on stopping x

emits langpicker_ready
emits langpicker_load_complete

env LOCALE_FILE=/var/local/system/locale
env ASR_TUTORIAL_COMPLETED_FILE=/var/local/system/asr_tutorial_launched
env FIRST_BOOT_FILE=/var/local/system/factory_fresh
env BLACKLIST=/opt/amazon/locale.blacklist

env LIBC_FATAL_STDERR_=1
env STDERR=/tmp/langpicker.err

task

script
    source /etc/upstart/env
    source /etc/upstart/functions

    if [ "$(devcap-get-feature -a frontlight)" -eq "1" ]; then
      lipc-set-prop com.lab126.powerd flStartup 1        
    fi

    if [ -e $BLACKLIST ]; then
        mntroot rw || true
        dsn=`cat /proc/usid | xargs echo -n`
        case $dsn in
            B05A*|B0F2* )
                # delete the locals specified in black list
                cat $BLACKLIST | xargs rm -rf
                RES=$?
                if [ $RES -ne 0 ]; then
                    f_log C locale prune "Failure while removing the blacklisted locales"
                else
                    f_log I locale prune "Successfully removed the blacklisted locales"
                    rm -rf $BLACKLIST
                fi
                ;;
           * )
                rm -rf $BLACKLIST
                ;;
        esac
        mntroot ro || true
    fi

    # see if ASR tutorial launched flag exists and 
    # delete the existing temp locale file already selected.
    if [ -e $ASR_TUTORIAL_COMPLETED_FILE ] && [ -e $LOCALE_FILE ] ; then
        rm -rf $LOCALE_FILE
    fi

    # see if we need to pick a language.
    # if the locale exists, then just exit
    if [ -e $LOCALE_FILE ]; then
        LANG=`awk -F'=' '/LANG/{print $2}' $LOCALE_FILE`
        LC_ALL=`awk -F'=' '/LC_ALL/{print $2}' $LOCALE_FILE`
        f_log I locale read "lang=$LANG,lc_all=$LC_ALL" "Retrieved Language"
        f_emit langpicker_ready LANG_PICKER_SCREEN_SHOWN=0
        exit 0
    fi

    # Did not find first boot file. This means that the device is booting in the factory line 
    # before the device has ever entered shipping mode. The lang picker should not be shown in this case.
    # Just set the locale to en-US.
    if [ ! -e $FIRST_BOOT_FILE ]; then
        # send the event to langpicker module to install the language
	lipc-send-event com.lab126.blanket.langpicker changeLocale -s "en-US"
	# wait for the localeChange event
	lipc-wait-event com.lab126.locale localeChange
	f_emit langpicker_ready LANG_PICKER_SCREEN_SHOWN=0
	exit 0
    fi

    # N.B. - only if we seem to have a valid input device.
    # the last input device in /etc/xorg.conf.* is the pointer, which is what we need.
    # we get "/dev/input/event3" - the extra eval is to strip quotes
    set $(grep /dev/input /etc/xorg.conf)
    INPUT=$(eval echo $(eval echo \$$#))

    if [ -e $INPUT -o -e /mnt/base-us/ENABLE_VNC ]; then
        # Start the pillow
        # TODO What if pillow fails to start ?        
        # setlocale(LC_ALL, "") depends on environment variable LC_ALL.
        export LC_ALL=en_US.UTF-8
        su framework -c "exec pillowd -NoDefaultPillowCase 2>> $STDERR"
        #Wait till pillow is ready before setting lipc property
        lipc-wait-event -s 30 com.lab126.pillow pillowReady
        # send the pick, wait for the change
        f_log I locale pick "" "Picking Language"
        lipc-set-prop com.lab126.pillow applicationWindow '{"name":"pickLocale","clientParams":{"show":true}}'
        f_emit langpicker_load_complete
        
        f_log I locale wait
        lipc-wait-event com.lab126.locale localeChange
        # Kill the pillow here, so that it boots up again with proper language.
        killall -q -s KILL pillowd || true
        f_log I locale set "" "User Picked Language"
    else
        f_log I locale no_input "" "no input device, defaulting"
        echo -e "LANG=en_US.UTF-8\nLC_ALL=en_US.UTF-8" > $LOCALE_FILE
    fi

    LANG=`awk -F'=' '/LANG/{print $2}' $LOCALE_FILE`
    LC_ALL=`awk -F'=' '/LC_ALL/{print $2}' $LOCALE_FILE`
    f_log I locale read "lang=$LANG,lc_all=$LC_ALL" "Set Language"
    f_emit langpicker_ready LANG_PICKER_SCREEN_SHOWN=1
end script

post-stop script
    source /etc/upstart/functions
    libc_scan
end script
