# rootfs.img.gz/fwo_rootfs.img/etc/upstart/touch.conf
start on mounted_varLocal and loaded_modules_dependencies

emits loaded_touch

task

script
 source /etc/upstart/functions
 /etc/upstart/touch && f_emit loaded_touch
end script
