# rootfs.img.gz/fwo_rootfs.img/etc/upstart/wandwrapper.conf
start on started lab126

script

  source /etc/upstart/functions
  f_wan
  start wand

end script
