# rootfs.img.gz/fwo_rootfs.img/etc/udev/rules.d/10-hw-wakeup.rules
# Hall Notifications
KERNEL=="rex_hall", DRIVER=="hall_sensor", RUN+="/lib/udev/bin/notifyUserTrigger"
# Powerbutton Notifications
KERNEL=="0-004b", DRIVER=="bd71827", RUN+="/lib/udev/bin/notifyUserTrigger"
