# rootfs.img.gz/fwo_rootfs.img/etc/udev/rules.d/20-waveform.rules
KERNEL=="imx_epdc_fb", SUBSYSTEM=="firmware", ACTION=="add", ATTR{loading}=="0", RUN+="/lib/udev/bin/waveform.sh", OPTIONS+="last_rule"
