# rootfs.img.gz/fwo_rootfs.img/etc/udev/rules.d/60-powerdtest.rules
# This rule is added for testing powerd
SUBSYSTEM=="platform", DRIVER=="max77696-battery", RUN+="/lib/udev/bin/notifysoc $ENV{DRIVER}"
